# backup_raid_local.timer
[Unit]
Description=Run RAID backup nightly

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=true

[Install]
WantedBy=timers.target
