# raid_backups.service
[Unit]
Description=Backup live RAID to backup disks

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/bin/local/backup_raid_local.sh
